# test-service/vigilix-test.service
[Unit]
Description=Vigilix Test Service

[Service]
ExecStart=/home/batman/.gemini/antigravity/scratch/vigilix/test-service/test-service.sh
Restart=always

[Install]
WantedBy=default.target
